# start_can.service
# /etc/systemd/system/my-startup.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/enable_CAN.sh

[Install]
WantedBy=multi-user.target
